# recipes-cizies-looper/cizies-looper/files/cizies-led-driver.service
[Unit]
Description=Cizies LED Driver startup

[Service]
ExecStart=/usr/sbin/cizies-led-driver

[Install]
WantedBy=multi-user.target
